# docker/supervisor/node_exporter.conf
[program:node_exporter]
; Use the 'server' command to start the metrics listener
; Point the scrape-uri to your pm.status_listen address (e.g., 127.0.0.1:9001)
command=/usr/local/bin/node_exporter --collector.disable-defaults --collector.cpu --collector.diskstats --collector.filesystem --collector.meminfo --collector.netstat

; Run as a non-privileged user
user=root
autostart=true
autorestart=true

; Adjust startsecs to ensure it stays up before considered 'running'
startsecs=5
startretries=3

; Logging setup
redirect_stderr=true
stdout_logfile=/var/log/node_exporter.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3