# docker/supervisor/php-fpm_exporter.conf
[program:php-fpm-exporter]
; Use the 'server' command to start the metrics listener
; Point the scrape-uri to your pm.status_listen address (e.g., 127.0.0.1:9001)
command=/usr/local/bin/php-fpm_exporter server --phpfpm.scrape-uri "unix:///run/php/php-fpm-status.sock;/fpm-status"

; Run as a non-privileged user
user=www-data
autostart=true
autorestart=true

; Adjust startsecs to ensure it stays up before considered 'running'
startsecs=5
startretries=3

; Logging setup
redirect_stderr=true
stdout_logfile=/var/log/php-fpm_exporter.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3